# college_de_montreal_ratings/college-ratings.service
[Unit]
Description=College de Montreal Teacher Ratings
After=network.target

[Service]
Type=simple
User=devstral
WorkingDirectory=/home/Devstral/Documents/chess/college_de_montreal_ratings
ExecStart=/usr/bin/python3 /home/Devstral/Documents/chess/college_de_montreal_ratings/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
